# dot_config/systemd/user/duckdns.service.tmpl
[Unit]
Description=Duck DNS Updater

[Service]
Type=simple
ExecStart={{ .chezmoi.homeDir }}/.local/bin/duckdns

[Install]
WantedBy=default.target
