# dot_config/systemd/user/foks.service
# Autogenerated by /usr/bin/foks; if you make custom modifications, delete this line to prevent overwrite.
[Unit]
Description=Foks background agent


[Service]
ExecStart=/usr/bin/foks agent
RestartSec=50ms
Restart=always

# Delays for 5s is SERVICE_RESULT != success
ExecStopPost=/usr/bin/foks ctl delay 5s

[Install]
WantedBy=default.target
