# dot_config/systemd/user/sworkstyle.service
[Unit]
Description=Renames workspaces to the icons configured

[Service]
ExecStart=/usr/bin/sworkstyle -d -l debug

