# dot_config/profile.d/executable_06-opencv.sh
# Disable warnings by OpenCV
export OPENCV_LOG_LEVEL=ERROR

