# dot_config/profile.d/executable_06-qlty-path.sh
# qlty
export QLTY_INSTALL="$HOME/.qlty"
export PATH="$QLTY_INSTALL/bin:$PATH"

