# dot_config/profile.d/executable_00-history-size.sh
export HISTSIZE=100000
export HISTFILESIZE=200000
# Zsh .zhistory size
export SAVEHIST=$HISTFILESIZE
