# dot_config/profile.d/executable_01-osx-manpath.sh
export MANPATH="$HOME/.rvm/man:/opt/X11/share/man:/usr/local/MacGPG2/share/man:/Library/Developer/CommandLineTools/usr/share/man:/usr/local/share/man:/usr/share/man:/usr/man"

