# dot_config/profile.d/executable_02-osx-gnu-screen-debug-path.sh
if [ -e "/opt/gnu-screen-debug/bin" ]; then
  case ":${PATH}:" in
    *:"/opt/gnu-screen-debug/bin":*) : ;; # no-op
    *) export PATH="/opt/gnu-screen-debug/bin:$PATH" ;;
  esac
fi

