# dot_config/profile.d/executable_01-editor-nvim.sh
# make default editor Neovim
#export EDITOR=nvim
## TODO: Migrate vim configs to Neovim...
## Until then... just use vim
export EDITOR=vim

