# dot_config/profile.d/executable_08-texlive.sh
export TEXMFDIST=/usr/share/texmf-dist
export TEXMFLOCAL=/usr/local/share/texmf:/usr/share/texmf
