# dot_config/profile.d/executable_00-xdg-terminal-command.sh
# set default terminal
export TERMINAL_COMMAND=xdg-terminal-exec

