# dot_config/profile.d/executable_01-osx-sbin-path.sh
case ":${PATH}:" in
  *:"/usr/local/sbin":*) : ;; # no-op
  *) export PATH="/usr/local/sbin:$PATH" ;;
esac

