# dot_config/profile.d/executable_07-ag-conf.sh
export AGRC_PATH="${XDG_CONFIG_HOME:-$HOME/.config}/ag/agrc"

