# dot_config/profile.d/executable_08-tpnote.sh.tmpl
# shellcheck shell=sh
# shellcheck disable=SC2034
export TPNOTE_USER='{{ .name }}'
