# dot_config/profile.d/executable_08-zeit.sh
# add default location for zeit.db
export ZEIT_DB="$HOME/.config/zeit.db"
