# dot_config/profile.d/executable_05-golang-path.sh
export GOPATH="$HOME/src/pub/go"
export PATH="$GOPATH/bin:$PATH"

