# dot_config/systemd/user/sway.service.d/00-sway-scripts-path.conf
[Service]
ExecStartPre=/usr/bin/systemctl --user set-environment "PATH=${HOME}/.local/share/sway/scripts:${HOME}/.local/bin:/usr/local/bin:/usr/bin"
